# gpio-control.service
[Unit]
Description=TonPI GPIO Control Service
After=network.target

[Service]
Type=simple
User={{SERVICE_USER}}
WorkingDirectory={{PROJECT_DIR}}
EnvironmentFile=-{{PROJECT_DIR}}/.env
SupplementaryGroups=spi gpio audio
ExecStart={{PHP_PATH}} artisan gpio:listen-controls
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target