# queue-worker.service
[Unit]
Description=TonPI Queue Worker
After=network.target
Wants=network-online.target

[Service]
Type=simple
User={{SERVICE_USER}}
WorkingDirectory={{PROJECT_DIR}}
ExecStart={{PHP_PATH}} artisan queue:work database --sleep=3 --tries=1 --timeout=30
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tonpi-queue
SupplementaryGroups=audio www-data

# Environment variables
Environment="PULSE_COOKIE=/var/run/pulse/cookie"
Environment="XDG_RUNTIME_DIR=/run/user/1000"

[Install]
WantedBy=multi-user.target
