# rfid-reader.service
[Unit]
Description=TonPI RFID Reader Service
After=network.target

[Service]
Type=simple
User={{SERVICE_USER}}
WorkingDirectory={{PROJECT_DIR}}
EnvironmentFile=-{{PROJECT_DIR}}/.env
SupplementaryGroups=spi gpio audio
ExecStart={{PHP_PATH}} artisan rfid:listen
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target